# crontabs/root
3,8,13,18,23,28,33,38,43,48,53,58 * * * * sleep 12 ; /bin/wget --no-check-certificate -O - http://freedns.afraid.org/dynamic/update.php?${YOUR_KEY} >> /tmp/freedns.log 2>&1 &
